# autoshutdown.service
[Unit]
Description = Automatically shutdown and destroy the server

[Service]
Type=simple
User=core
ExecStart=/tmp/chron.sh
